# deploy/macos/com.lex.splitbot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.lex.splitbot</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Users/alexanderroller/BurryDeez/scripts/start_splitbot.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>WorkingDirectory</key>
    <string>/Users/alexanderroller/BurryDeez</string>

    <key>StandardOutPath</key>
    <string>/Users/alexanderroller/BurryDeez/logs/splitbot.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/alexanderroller/BurryDeez/logs/splitbot.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PYTHONUNBUFFERED</key>
      <string>1</string>
    </dict>
  </dict>
</plist>
